# api/crontab
# Выполнение каждый час с 10 до 18
0 10-18 * * * /app/script.sh >> /app/logfile.log 2>&1
