# api/supervisord.conf
[supervisord]
nodaemon=true

[program:my_script]
command=/app/script.sh
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/my_script.log

[program:cron]
command=cron -f -L /dev/stdout