# api/supervisord.conf
[supervisord]
nodaemon=true

[program:script]
command=/bin/sh -c "while true; do /bin/sh /app/script.sh; sleep 55m; done"

[program:uvicorn]
command=uvicorn main:app --host 0.0.0.0 --port 8000